# henlo.service
[Unit]
Description=Henlo backend service
After=network.target

[Service]
ExecStart=/usr/bin/docker stack deploy --compose-file=docker-compose.yml henlo
WorkingDirectory=/home/pi/henlo-backend
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
